# app/supervisor.conf
[program:app]
command=python3.6 gserver.py
directory=/home/deploy/catyy/app
user=root