# app/supervisor.conf
[program:webapp]
command=uwsgi uwsgi.ini
directory=/home/deploy/catyy/app
user=root